# labs/oscp/machines/windows/easy/legacy/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

# SSH Service
[program:ssh]
command=/usr/sbin/sshd -D
stdout_logfile=/var/log/supervisor/ssh.log
stderr_logfile=/var/log/supervisor/ssh.log
autorestart=true
priority=100

# Legacy Services Startup
[program:legacy-services]
command=/opt/legacy/start-services.sh
stdout_logfile=/var/log/supervisor/legacy-services.log
stderr_logfile=/var/log/supervisor/legacy-services.log
autorestart=true
priority=200
startsecs=10

# Samba SMB Service
[program:smbd]
command=/usr/sbin/smbd --foreground --no-process-group
stdout_logfile=/var/log/supervisor/smbd.log
stderr_logfile=/var/log/supervisor/smbd.log
autorestart=true
priority=300
startsecs=5

# Samba NetBIOS Service
[program:nmbd]
command=/usr/sbin/nmbd --foreground --no-process-group
stdout_logfile=/var/log/supervisor/nmbd.log
stderr_logfile=/var/log/supervisor/nmbd.log
autorestart=true
priority=300
startsecs=5

# MS17-010 Vulnerable SMB Simulation
[program:vulnerable-smb]
command=python3 /opt/legacy/vulnerable-smb.py
directory=/opt/legacy
stdout_logfile=/var/log/supervisor/vulnerable-smb.log
stderr_logfile=/var/log/supervisor/vulnerable-smb.log
autorestart=true
priority=400
startsecs=10

# System Monitor (for educational purposes)
[program:system-monitor]
command=/bin/bash -c "while true; do echo '[$(date)] Legacy machine running - IP: $(hostname -I)'; sleep 300; done"
stdout_logfile=/var/log/supervisor/system-monitor.log
stderr_logfile=/var/log/supervisor/system-monitor.log
autorestart=true
priority=900